# ops/ourdigest/ourdigest-refresh.service
[Unit]
Description=ourdigest refresh - fetch sources and write feeds
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=architect
WorkingDirectory=/home/architect/ourdigest
EnvironmentFile=/home/architect/ourdigest/.env
ExecStart=/home/architect/ourdigest/.venv/bin/ourdigest fetch
Nice=10
